# autofs/mnt-synology.automount
[Unit]
Description=Automount Nas Synology

[Automount]
Where=/mnt/synology

[Install]
WantedBy=multi-user.target
